# worldeditor-pkg/src/deb/control/postrm
#!/bin/sh
 	
update-rc.d worldeditor remove
exit 0